# firstmail.service
[Unit]
Description=FirstMail Tool (批量改密码 / 邮件查看器)
After=network.target

[Service]
Type=simple
# 部署目录，按实际修改
WorkingDirectory=/opt/firstmail
# node 路径用 `which node` 确认；nodesource 安装一般是 /usr/bin/node
ExecStart=/usr/bin/node server.js
Restart=always
RestartSec=3
# 运行用户，按需改成你自己的用户（如 ubuntu）
User=www-data
Group=www-data

# 配置方式二选一：
#   A) 在部署目录放 .env（推荐），server.js 会自动读取
#   B) 取消下面注释，直接在这里注入环境变量
# Environment=FIRSTMAIL_API_KEY=your_api_key_here
# Environment=PORT=8787
# Environment=HOST=127.0.0.1

# 安全加固（可选）
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
